# vagrant-config/worker.conf
[program:worker]
command=php artisan queue:listen
directory=/var/www/
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/swibeco-backend-log.txt
